# deployment/supervisor-worker.conf
# Image Editor - Supervisor Configuration
# Copyright © 2026 Md. Sahabul. All rights reserved.
# Designed & developed by Md. Sahabul.
#
# Place at /etc/supervisor/conf.d/image-editor.conf
# Then: sudo supervisorctl reread && sudo supervisorctl update

; ── Laravel queue workers ─────────────────────────────────────────────
[program:image-editor-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/image-editor/backend/artisan queue:work redis \
    --sleep=3 \
    --tries=3 \
    --timeout=180 \
    --max-jobs=500 \
    --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/image-editor-worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stopwaitsecs=300

; ── AI inference microservice ─────────────────────────────────────────
[program:image-editor-ai]
command=/var/www/image-editor/ai-service/venv/bin/uvicorn main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1 \
    --timeout-keep-alive 120
directory=/var/www/image-editor/ai-service
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/image-editor-ai.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
